# ops/configs/liquidsoap-radio.service
[Unit]
Description=Liquidsoap AutoDJ (Pipe to FFmpeg)
After=network.target

[Service]
# The Magic: Liquidsoap | FFmpeg -> RTMP
ExecStart=/bin/bash -c '/usr/bin/liquidsoap /etc/liquidsoap/autodj.liq | /usr/bin/ffmpeg -f wav -i pipe:0 -c:a aac -b:a 128k -f flv rtmp://127.0.0.1/autodj/stream'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
